# daemon/launchd/com.somniatore.claude-relay.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.somniatore.claude-relay</string>
  <key>ProgramArguments</key>
  <array>
    <string>/Users/ghaz/.nvm/versions/node/v22.22.0/bin/node</string>
    <string>/Users/ghaz/projects/startups/claude-relay/daemon/dist/index.js</string>
  </array>
  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <true/>
  <key>StandardOutPath</key>
  <string>/Users/ghaz/Library/Logs/claude-relay/daemon.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/ghaz/Library/Logs/claude-relay/daemon.log</string>
  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
  </dict>
</dict>
</plist>
